# ansible/roles/backup/templates/docker-secure-platform-backup.service.j2
[Unit]
Description=Docker Secure Platform backup

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/docker-secure-platform/scripts/backup.sh
